# recipes/rclone/mnt-rclone.automount
[Unit]
After=remote-fs.target

[Automount]
Where=/mnt/rclone
TimeoutIdleSec=600

[Install]
WantedBy=remote-fs.target
